# hourroi-bot.service
[Unit]
Description=HOURROI Telegram Bot - Blockchain Event Monitor
Documentation=https://github.com/your-username/hourroi-bot
After=network.target
Wants=network.target

[Service]
Type=forking
PIDFile=/opt/hourroi-bot/.pm2/pm2.pid
User=www-data
Group=www-data

# PM2 commands
ExecStart=/usr/bin/pm2 start /opt/hourroi-bot/ecosystem.config.js --env production
ExecReload=/usr/bin/pm2 reload all
ExecStop=/usr/bin/pm2 kill

# Working directory
WorkingDirectory=/opt/hourroi-bot

# Environment
Environment=NODE_ENV=production
Environment=PATH=/usr/bin:/usr/local/bin
Environment=PM2_HOME=/opt/hourroi-bot/.pm2

# Restart settings
Restart=always
RestartSec=10
StartLimitInterval=60s
StartLimitBurst=3

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/hourroi-bot /var/log/hourroi-bot /tmp

# Resource limits
LimitNOFILE=65536
LimitNPROC=32768

[Install]
WantedBy=multi-user.target